# src/etc/systemd/system/munin-node.service.d/override.conf
[Unit]
After=
After=vyos-router.service

[Service]
ExecStart=""
ExecStart="${VRF_ARGS} /usr/sbin/munin-node --foreground $DAEMON_ARGS"
ExecReload=/bin/kill -HUP $MAINPID
